# deploy/exec/crontab
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.XXXX1pk9fJ installed on Thu Mar 13 12:32:30 2014)
# (Cron version V5.0 -- $Id: crontab.c,v 1.12 2004/01/23 18:56:42 vixie Exp $)
1 * * * * sh /var/www/poker.sch.bme.hu/root/exec/crontab_testbots.sh &> /dev/null
10 * * * * sh /var/www/poker.sch.bme.hu/root/exec/crontab_update.sh &> /dev/null
